# etc/systemd/system/storm-webdav.service.d/storm-webdav.conf
[Service]
# Storm webdav service will run as this user
# Default: storm
Environment="STORM_WEBDAV_USER=storm"

# JVM options
Environment="STORM_WEBDAV_JVM_OPTS=-Xms1024m -Xmx1024m"
# To have the JVM debug agent bind on all interfaces use `...,address=*:1044,suspend=n`, otherwise the agent binds only on
# localhost
# Environment="STORM_WEBDAV_JVM_OPTS=-Xms1024m -Xmx1024m -Xdebug -Xrunjdwp:server=y,transport=dt_socket,address=1044,suspend=n"

# StoRM webdav will bind on this address
Environment="STORM_WEBDAV_SERVER_ADDRESS=0.0.0.0"

# In order to handle tpc requests, webdav needs to know for which hostnames
# is handling requests (and possibly any alias)
# STORM_WEBDAV_HOSTNAME_0="localhost"
# STORM_WEBDAV_HOSTNAME_1="alias.for.localhost"

# HTTPS connector port
Environment="STORM_WEBDAV_HTTPS_PORT=8443"

# HTTP connector port
Environment="STORM_WEBDAV_HTTP_PORT=8085"

# Path to the service certificate.
Environment="STORM_WEBDAV_CERTIFICATE_PATH=/etc/grid-security/storm-webdav/hostcert.pem"

# Path to the service certificate private key.
Environment="STORM_WEBDAV_PRIVATE_KEY_PATH=/etc/grid-security/storm-webdav/hostkey.pem"

# Trust anchors directory
Environment="STORM_WEBDAV_TRUST_ANCHORS_DIR=/etc/grid-security/certificates"

# How often (seconds) trust anchors are refreshed from the trust anchors directory
Environment="STORM_WEBDAV_TRUST_ANCHORS_REFRESH_INTERVAL=86400"

# Max concurrent connections
Environment="STORM_WEBDAV_MAX_CONNECTIONS=300"

# Connection queue size
Environment="STORM_WEBDAV_MAX_QUEUE_SIZE=900"

# Connector Maximum idle time (in milliseconds)
Environment="STORM_WEBDAV_CONNECTOR_MAX_IDLE_TIME=30000"

# Storage area configuration directory
Environment="STORM_WEBDAV_SA_CONFIG_DIR=/etc/storm/webdav/sa.d"

# The webdav server jar. Do not change this unless you know what you are doing.
Environment="STORM_WEBDAV_JAR=/usr/share/java/storm-webdav/storm-webdav-server.jar"

# The StoRM webdav logfile. Do not change this unless you know what you are doing.
Environment="STORM_WEBDAV_LOG=/var/log/storm/webdav/storm-webdav-server.log"

# The StoRM webdav stdout logfile. Do not change this unless you know what you are doing.
Environment="STORM_WEBDAV_OUT=/var/log/storm/webdav/storm-webdav-server.out"

# The StoRM webdav stderr logfile. Do not change this unless you know what you are doing.
Environment="STORM_WEBDAV_ERR=/var/log/storm/webdav/storm-webdav-server.err"

# The StoRM webdav logging configuration. Do not change this unless you know what you are doing.
Environment="STORM_WEBDAV_LOG_CONFIGURATION=/etc/storm/webdav/logback.xml"

# The StoRM webdav access log configuration. Do not change this unless you know what you are doing.
Environment="STORM_WEBDAV_ACCESS_LOG_CONFIGURATION=/etc/storm/webdav/logback-access.xml"

# Enables/disables  support for VO map files
Environment="STORM_WEBDAV_VO_MAP_FILES_ENABLE=false"

# Where StoRM WebDAV should look for voms map files.
Environment="STORM_WEBDAV_VO_MAP_FILES_CONFIG_DIR=/etc/storm/webdav/vo-mapfiles.d"

# VOMS map files refresh period (in seconds)
Environment="STORM_WEBDAV_VO_MAP_FILES_REFRESH_INTERVAL=21600"

# Max number of active third party copy connections
Environment="STORM_WEBDAV_TPC_MAX_CONNECTIONS=50"

# Proactively verify checksums in third party copy
Environment="STORM_WEBDAV_TPC_VERIFY_CHECKSUM=false"

# Max number of active third party copy connections per route
Environment="STORM_WEBDAV_TPC_MAX_CONNECTIONS_PER_ROUTE=25"

# Enables the OAuth authorization server, used to issue bearer tokens
# to clients presenting a VOMS credential
# Environment="STORM_WEBDAV_AUTHZ_SERVER_ENABLE=false"

# The issuer setting for the token
# Default: https://${STORM_WEBDAV_HOSTNAME_0}:${STORM_WEBDAV_HTTPS_PORT}
# Environment="STORM_WEBDAV_AUTHZ_SERVER_ISSUER=https://storm.example:8443"

# Maximum token lifetime in seconds
# Environment="STORM_WEBDAV_AUTHZ_SERVER_MAX_TOKEN_LIFETIME_SEC=43200"

# The WebDAV authz server secret
# This secret is used to sign bearer tokens, so if you have multiple
# instances of storm-webdav serving for the same VOs, you should use
# the same secret
# IMPORTANT: set this in production, do not rely on the default random
# value as that will cause problems across service restarts
#
# Environment="STORM_WEBDAV_AUTHZ_SERVER_SECRET=areasonablesecretlongerthan256bits"

# Should StoRM WebDAV always require a valid client certificate on the HTTPS endpoint?
# Set to 'false' if you want token-based authentication to work (and thus third-party copy)
# Environment="STORM_WEBDAV_REQUIRE_CLIENT_CERT=false"

# Use Conscrypt JSSE provider?
# Default: false
# Set to 'true' if you want to use a faster TLS implementation which can improve dramatically
# the performance on data transfers. This flag must also set to "true" when HTTP/2 is enabled
# Environment="STORM_WEBDAV_USE_CONSCRYPT=false"

# Use Conscrypt JSSE provider for TPC client connections?
# Default: false
# Set to 'true' if you want to use a faster TLS implementation which can improve dramatically
# the performance on data transfers. This flag controls the use of Conscrypt for TPC client connections.
# Environment="STORM_WEBDAV_TPC_USE_CONSCRYPT=false"

# Enable HTTP/2 support?
# Default: false
# Set to 'true' if you want to enable HTTP/2 (and remember to enable conscrypt too!)
# Environment="STORM_WEBDAV_ENABLE_HTTP2=false"

# Source file for the tape REST API well-known endpoint
# Default: '/etc/storm/webdav/wlcg-tape-rest-api.json'
# Environment="STORM_WEBDAV_TAPE_WELLKNOWN_SOURCE=/etc/storm/webdav/wlcg-tape-rest-api.json"

# Buffer size for both internal and third-party copy requests.
# This adds more efficiency than to write the whole data. Valid values are numbers >= 4096.
# Default: 1048576
# Environment="STORM_WEBDAV_BUFFER_FILE_BUFFER_SIZE_BYTES=1048576"

# Enable checksum filter which adds checksum as an header following RFC 3230.
# Default: true
# Environment="STORM_WEBDAV_CHECKSUM_FILTER_ENABLED=true"

# Enable Macaroon filter to process Macaroon tokens. Requires authz server enabled.
# Default: true
# Environment="STORM_WEBDAV_MACAROON_FILTER_ENABLED=true"

# TLS protocol for non-TPC requests
# Default: TLS
# Environment="STORM_WEBDAV_TLS_PROTOCOL=TLS"

# VOMS Trust Store directory
# Default: /etc/grid-security/vomsdir
# Environment="STORM_WEBDAV_VOMS_TRUST_STORE_DIR=/etc/grid-security/vomsdir"

# VOMS Trust Store refresh interval
# Default: 43200
# Environment="STORM_WEBDAV_VOMS_TRUST_STORE_REFRESH_INTERVAL_SEC=43200"

# Enable caching for VOMS certificate validation
# Default: true
# Environment="STORM_WEBDAV_VOMS_CACHE_ENABLE=true"

# Cache entries lifetime, used if caching for VOMS certificate validation is enabled
# Default: 300
# Environment="STORM_WEBDAV_VOMS_CACHE_ENTRY_LIFETIME_SEC=300"

# Enable SciTags support
# Default: false
# Environment="STORM_WEBDAV_SCITAGS_ENABLED=false"
